# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/bin/bash -c "/usr/bin/ssh-keygen -A && /usr/sbin/sshd -D"
